# docker/django/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/app/logs/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/tmp/supervisord.pid

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[program:gunicorn]
command=gunicorn money_journey.wsgi:application --bind 0.0.0.0:8000 --workers 3 --access-logfile /app/logs/gunicorn_access.log --error-logfile /app/logs/gunicorn_error.log
directory=/app
user=django
autostart=true
autorestart=true
redirect_stderr=false
stdout_logfile=/app/logs/gunicorn.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/app/logs/gunicorn_error.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
stopasgroup=true
killasgroup=true
environment=LANG=en_US.UTF-8,LC_ALL=en_US.UTF-8,PATH="/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

[program:cron]
command=/usr/sbin/cron -f
user=root
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/app/logs/cron.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stopasgroup=true
killasgroup=true
environment=PATH="/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

[group:django]
programs=gunicorn,cron
priority=999
